# judge/create-cgroups.service.in
[Unit]
Description=Make sure cgroups exist for DOMjudge judgedaemon

[Service]
Type=oneshot
ExecStart=@judgehost_bindir@/create_cgroups
RemainAfterExit=true
