# judge/domjudge-judgehost.service.in
[Unit]
Description=DOMjudge JudgeDaemon
Requires=create-cgroups.service
After=create-cgroups.service
After=network.target

[Service]
Type=simple

ExecStart=@judgehost_bindir@/judgedaemon -n 0
User=@DOMJUDGE_USER@

Restart=always
RestartSec=3
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
