# judge/domjudge-judgehost.service.in
[Unit]
Description=DOMjudge JudgeDaemon
Requires=create-cgroups.service
After=create-cgroups.service
After=network.target

[Service]
Type=simple

ExecStart=@judgehost_bindir@/judgedaemon
User=@DOMJUDGE_USER@

# If you want to run multiple judgedaemons per machine, or bind it to
# a single CPU core, use the following instead, replacing X with the
# specific core, starting from 0 and create a separate copy of this
# file for each judgedaemon instance.
#ExecStart=@judgehost_bindir@/judgedaemon -n X
#User=@DOMJUDGE_USER@-X

Restart=always
RestartSec=3
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
