# gitlab/demodocker/supervisor/nginx.conf
[program:nginx]
command=nginx -g "daemon off;"
numprocs=1
autostart=true
autorestart=true
