# gitlab/demodocker/supervisor/php.conf
[program:php]
command=php-fpm -F
numprocs=1
autostart=true
autorestart=true
