# etc/domjudge-judgehost.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          domjudge-judgehost
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $all $network apache2 nginx mysql
# Should-Stop:       $all $network apache2 nginx mysql
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: DOMjudge Judgedaemon
### END INIT INFO
#
# SystemV init script to start a judgedaemon.
# In OPTIONS you can configure which CPUset to use when using cgroups.
#
# This is currently very much geared towards Debian(like) systems.
# Improvements for other systems welcome.

NAME=judgedaemon
OPTIONS='-d -v3 -n 0'

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
BINDIR=@judgehost_bindir@
RUNDIR=@judgehost_rundir@
CHROOTDIR=@judgehost_chrootdir@
CGROUPDIR=@judgehost_cgroupdir@
DAEMON=$BINDIR/judgedaemon
PIDFILE=$RUNDIR/$NAME.pid
DESC="DOMjudge Judgedaemon"
RUNAS=@DOMJUDGE_USER@

test -x $DAEMON || exit 0

# shellcheck disable=SC1091
. /lib/lsb/init-functions

if [ ! -d $RUNDIR ]; then
	mkdir -p $RUNDIR
	chown $RUNAS $RUNDIR
fi

# Make sure cgroups are mounted where we expect them
if [ ! -d $CGROUPDIR/memory/domjudge ]; then
	$BINDIR/create_cgroups
fi

# Check for required cgroup feature
if [ ! -f $CGROUPDIR/memory/domjudge/memory.memsw.limit_in_bytes ]; then
	echo "Error: Swapaccounting support missing from cgroups. Unable to continue.

To fix this, please make the following changes:
1. In /etc/default/grub, add 'swapaccount=1' to GRUB_CMDLINE_LINUX_DEFAULT
2. Run update-grub
3. Reboot" >&2
exit 0
fi

# Check that the chroot is created/set up properly
if [ ! -d $CHROOTDIR ]; then
	echo "Warning: chroot $CHROOTDIR not created. Run dj_make_chroot to create it."
fi


# Close some file descriptors beyond stdio, because these cannot be
# closed in the judgedaemon PHP script, workaround for Debian postinst
# hanging on open fd.
exec 3>&-
exec 4>&-
exec 5>&-

case "$1" in
  start)
	log_begin_msg "Starting $DESC: $NAME"
	if [ -d $CHROOTDIR ]; then
		mount -t proc proc $CHROOTDIR/proc
	fi
	start-stop-daemon --start --quiet --chuid $RUNAS --pidfile "$PIDFILE" --exec "$DAEMON" -- $OPTIONS
	log_end_msg $?
	;;
  stop)
	log_begin_msg "Stopping $DESC: $NAME"
	start-stop-daemon --stop --quiet --retry 5 --pidfile $PIDFILE
	if [ -d $CHROOTDIR ]; then
		umount $CHROOTDIR/proc || true
	fi
	log_end_msg $?
	;;
  reload|force-reload|restart)
	log_begin_msg "Restarting $DESC: $NAME"
	if start-stop-daemon --stop --quiet --retry 5 --pidfile $PIDFILE; then
		start-stop-daemon --start --quiet --chuid $RUNAS --pidfile "$PIDFILE" --exec "$DAEMON" -- $OPTIONS
	fi
	log_end_msg $?
	;;
  status)
	# shellcheck disable=SC2059
	printf "Status of $DESC: "
	if [ ! -r "$PIDFILE" ]; then
		echo "$NAME is not running."
		exit 3
	fi
	if read -r pid < "$PIDFILE" && ps -p "$pid" > /dev/null 2>&1; then
		echo "$NAME is running."
		exit 0
	else
		echo "$NAME is not running but $PIDFILE exists."
		exit 1
	fi
	;;
  *)
	N=/etc/init.d/${0##*/}
	echo "Usage: $N {start|stop|restart|force-reload|status}" >&2
	exit 1
	;;
esac

exit 0
